# katago-server.service
[Unit]
Description=KataGo Server (Rust)
After=network.target

[Service]
Type=simple
User=your-user-name
Group=your-group-name
WorkingDirectory=/var/www/katago-server
ExecStart=/var/www/katago-server/target/release/katago-server
Restart=on-failure
RestartSec=10

# Environment
Environment="RUST_LOG=info"

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/katago-server

[Install]
WantedBy=multi-user.target
